# sources/akmods/akmods-keygen.target
[Unit]
Wants=akmods-keygen@.service
PartOf=akmods.service
