# sources/akmods/akmods.service.in
[Unit]
Description=Builds and install new kmods from akmod packages
ConditionPathExists=!/run/ostree-booted
Before=@SERVICE@
After=akmods-keygen.target
Wants=akmods-keygen.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/akmods --from-init

[Install]
WantedBy=multi-user.target
